# iso/airootfs/etc/systemd/system/amd.service
[Unit]
Description=Enable AMD mode

[Service]
ExecStart=/bin/bash /etc/systemd/system/amd.sh

[Install]
WantedBy=multi-user.target
